# daemons/supervisor.d/gunicorn.ini
[program:gunicorn]
command=python manage.py run_gunicorn -p ./daemons/tmp/gunicorn.pid -w 2
directory=../
autostart=true
autorestart=true
stdout_logfile=./logs/gunicorn.log
redirect_stderr=true
