# daemons/supervisor.d/ztaskd.ini
[program:ztaskd]
command=python manage.py ztaskd ztaskd --noreload -f daemons/logs/ztaskd.log
directory=../
autostart=true
autorestart=true
stdout_logfile=./logs/ztaskd.log
redirect_stderr=true
